# package/siteview/vppn-mqtt/init/mqtt-client.init
#!/bin/sh /etc/rc.common
START=99

start()
{
	vpnlogd
	mqtt-client 
	mqtt-client -m
}

stop()
{
	killall vpnlogd
	killall mqtt-client
}
