# package/siteview/vppn-mqtt/init/vppnctrl.init
#!/bin/sh /etc/rc.common
START=99

MAX_TUNNEL_NUM=1

start()
{
	#vppnconfig loadconfig vppn
	for i in $(seq ${MAX_TUNNEL_NUM}); do
		vppnctrl -t $(( i - 1 )) -c vppn run
	done
	#/usr/sbin/vpnmonitor
	#vpnctrl dns_report
}

stop()
{
	for i in $(seq ${MAX_TUNNEL_NUM}); do
		vppnctrl -t $(( i - 1 )) -c vppn kill
	done
	killall vpnmonitor
}
