# package/siteview/webcgi/files/etc/uci-defaults/01_sysdep-defaults
#!/bin/sh

[ -e /etc/config/network ] && exit 0

touch /etc/config/network

/sbin/sysinit network_uci_defaults

exit 0
