# package/siteview/webcgi/files/etc/uci-defaults/01_sysdep-defaults
#!/bin/sh

/sbin/sysinit network_uci_defaults

exit 0
