# package/siteview/webcgi/files/etc/hotplug.d/iface/99-route
#!/bin/sh

[ "$ACTION" = "ifup" ] && {
	/etc/init.d/route restart
}
